# files/nullmailer.cron
33 3 * * * root find /var/spool/nullmailer/queue/ -type f -mtime +1 -delete